# packaging/macos/com.canonical.multipassd.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.canonical.multipassd</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>@CPACK_PACKAGING_INSTALL_PREFIX@/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string>@CPACK_PACKAGING_INSTALL_PREFIX@/bin/multipassd</string>
      <string>--verbosity</string>
      <string>debug</string>
    </array>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Library/Logs/Multipass/multipassd.log</string>
    <key>StandardErrorPath</key>
    <string>/Library/Logs/Multipass/multipassd.log</string>

    <key>ExitTimeOut</key>
    <integer>45</integer>
</dict>

</plist>
